# ops/systemd/hateful-eight-refresh.timer
[Unit]
Description=Run Hateful Eight refresh shortly after US market close

[Timer]
OnCalendar=Mon..Fri 16:35:00 America/New_York
Persistent=true
RandomizedDelaySec=0
AccuracySec=1m

[Install]
WantedBy=timers.target
