# ops/systemd/hateful-eight-refresh.service
[Unit]
Description=Hateful Eight nightly data refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=paul
Group=paul
Environment=HOME=/home/paul
WorkingDirectory=/srv/repos/hateful-eight
ExecStart=/srv/repos/hateful-eight/ops/scripts/daily_refresh.sh

